# infra/m1/com.logan.goose-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.logan.goose-nightly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/namto/.local/bin/goose</string>
    <string>run</string>
    <string>--recipe</string>
    <string>/Users/namto/Desktop/agentic-platform/recipes/nightly-report.yaml</string>
    <string>--no-session</string>
    <string>--max-turns</string>
    <string>15</string>
    <string>--params</string>
    <string>REPOS=namlogan/agentic-platform</string>
    <string>--params</string>
    <string>PLATFORM_REPO=namlogan/agentic-platform</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/namto/Desktop/agentic-platform</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GOOSE_PROVIDER</key><string>ollama</string>
    <key>GOOSE_MODEL</key><string>qwen2.5:7b</string>
    <key>OLLAMA_HOST</key><string>http://milai:11434</string>
    <key>PATH</key><string>/Users/namto/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key><string>/Users/namto</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>RunAtLoad</key><false/>
  <key>KeepAlive</key><false/>
  <key>StandardOutPath</key><string>/Users/namto/Library/Logs/goose-nightly.log</string>
  <key>StandardErrorPath</key><string>/Users/namto/Library/Logs/goose-nightly.err</string>
</dict></plist>
