# infra/rtx/vllm.service
[Unit]
Description=vLLM OpenAI-compatible server
After=network-online.target docker.service
Wants=network-online.target

[Service]
User=vllm
Restart=always
RestartSec=10
ExecStart=/usr/bin/docker run --rm \
  --gpus all \
  --name vllm-server \
  -p 8000:8000 \
  -v /opt/models:/models \
  vllm/vllm-openai:latest \
    --model Qwen/Qwen3-32B-AWQ \
    --max-model-len 32768 \
    --gpu-memory-utilization 0.92 \
    --served-model-name Qwen/Qwen3-32B-AWQ \
    --trust-remote-code
ExecStop=/usr/bin/docker stop vllm-server
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
